# supervisord.conf
; ============================================================
; Supervisord 配置 — 知识库问答系统
; Profile: single-gpu (单张 RTX 5080 (16GB) 部署全部模型，适合开发测试环境)
; 部署模式: native
; 自动生成，请勿手动编辑
; ============================================================
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
logfile=/tmp/supervisord.log
logfile_maxbytes=50MB
logfile_backups=10
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:redis]
command=redis-server --appendonly yes --maxmemory 512mb --maxmemory-policy allkeys-lru
directory=/mnt/d/Work/AI/AGIC/aigic
environment=
priority=1
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/redis.log
stdout_logfile_maxbytes=50MB


[program:chromadb]
command=chroma run --host 0.0.0.0 --port 8004 --path ./data/chroma
directory=/mnt/d/Work/AI/AGIC/aigic
environment=IS_PERSISTENT="TRUE",ANONYMIZED_TELEMETRY="FALSE"
priority=2
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/chromadb.log
stdout_logfile_maxbytes=50MB


[program:llm-service]
command=bash services/llm-service/scripts/start_vllm.sh
directory=/mnt/d/Work/AI/AGIC/aigic
environment=CUDA_VISIBLE_DEVICES="0",MODEL_NAME="google/gemma-4-e2b-it",TENSOR_PARALLEL_SIZE="1",QUANTIZATION="int4",CONTEXT_WINDOW="2048",MAX_NUM_SEQS="1",GPU_MEMORY_UTILIZATION="0.75",KV_CACHE_DTYPE="auto",PORT="8001",MODEL_PATH="/mnt/d/Work/AI/AGIC/aigic/models/gemma3"
priority=10
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/llm-service.log
stdout_logfile_maxbytes=50MB


[program:embedding-service]
command=python services/embedding-service/app/main.py
directory=/mnt/d/Work/AI/AGIC/aigic
environment=CUDA_VISIBLE_DEVICES="0",MODEL_NAME="BAAI/bge-m3",MAX_LENGTH="8192",BATCH_SIZE="8",PORT="8002"
priority=11
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/embedding-service.log
stdout_logfile_maxbytes=50MB


[program:asr-service]
command=python services/asr-service/app/main.py
directory=/mnt/d/Work/AI/AGIC/aigic
environment=CUDA_VISIBLE_DEVICES="0",MODEL_NAME="large-v3",DEFAULT_LANGUAGE="zh",VAD_THRESHOLD="0.5",PORT="8005"
priority=12
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/asr-service.log
stdout_logfile_maxbytes=50MB


[program:tts-service]
command=python services/tts-service/app/main.py
directory=/mnt/d/Work/AI/AGIC/aigic
environment=CUDA_VISIBLE_DEVICES="0",DEFAULT_VOICE="default",SAMPLE_RATE="24000",PORT="8006"
priority=13
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/tts-service.log
stdout_logfile_maxbytes=50MB


[program:rag-service]
command=python services/rag-service/app/main.py
directory=/mnt/d/Work/AI/AGIC/aigic
environment=EMBEDDING_SERVICE_URL="http://localhost:8002",CHROMA_HOST="localhost",CHROMA_PORT="8004",CHROMA_COLLECTION="kb_knowledge",TOP_K="5",RERANK_TOP_N="3",SIMILARITY_THRESHOLD="0.65",REDIS_HOST="localhost",REDIS_PORT="6379",PORT="8003"
priority=20
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/rag-service.log
stdout_logfile_maxbytes=50MB


[program:api-gateway]
command=uvicorn services.api-gateway.app.main:app --host 0.0.0.0 --port 8000
directory=/mnt/d/Work/AI/AGIC/aigic
environment=LLM_SERVICE_URL="http://localhost:8001",EMBEDDING_SERVICE_URL="http://localhost:8002",RAG_SERVICE_URL="http://localhost:8003",ASR_SERVICE_URL="http://localhost:8005",TTS_SERVICE_URL="http://localhost:8006",REDIS_HOST="localhost",REDIS_PORT="6379",CONVERSATION_TTL="86400",PORT="8000"
priority=30
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/api-gateway.log
stdout_logfile_maxbytes=50MB


[program:frontend]
command=npm run dev
directory=/mnt/d/Work/AI/AGIC/aigic/frontend
environment=NEXT_PUBLIC_API_URL="http://localhost:8000",NEXT_PUBLIC_WS_URL="ws://localhost:8000",PORT="3000"
priority=40
autostart=true
autorestart=true
startsecs=5
startretries=3
stopwaitsecs=10
stopsignal=TERM
redirect_stderr=true
stdout_logfile=/tmp/frontend.log
stdout_logfile_maxbytes=50MB


